# src/Vagrant/vagrant-g07/Scripts_FilesMatrix/webserver_message.timer
[Unit]
Description=Run Matrix notification script at specific times

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true

[Install]
WantedBy=timers.target